# deploy/systemd/broadcast-monitor.service
[Unit]
Description=Broadcast Uptime Monitor
After=network.target

[Service]
Type=simple
User=travisr
WorkingDirectory=/opt/broadcast-uptime-monitor
Environment=PATH=/opt/broadcast-uptime-monitor/.venv/bin:/usr/bin:/bin
ExecStart=/opt/broadcast-uptime-monitor/.venv/bin/python /opt/broadcast-uptime-monitor/scripts/monitor.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
